# rootfs/etc/services.d/mordhao/run
#!/usr/bin/execlineb -P
/mnt/steam/mordhau/MordhauServer.sh
