# rootfs/etc/services.d/mordhao/run
#!/usr/bin/execlineb -P
tail -f /dev/null
